# deploy/autodl/supervisord.conf
[unix_http_server]
file=runtime/supervisor/supervisor.sock
chmod=0700

[supervisord]
logfile=logs/supervisord.log
pidfile=runtime/supervisor/supervisord.pid
childlogdir=logs
loglevel=info
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://runtime/supervisor/supervisor.sock

[program:api]
directory=.
command=.venv/bin/python -m uvicorn gateway.main:app --host 127.0.0.1 --port 8012
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=10
stopasgroup=true
killasgroup=true
stdout_logfile=logs/api.out.log
stderr_logfile=logs/api.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:quant-scheduler]
directory=.
command=.venv/bin/python scripts/quant_signal_scheduler.py --poll-seconds 20
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=10
startretries=10
stopasgroup=true
killasgroup=true
stdout_logfile=logs/scheduler.out.log
stderr_logfile=logs/scheduler.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
